# crontab
PATH="/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
*/5 * * * * sh /opt/Industriallog/run.sh
